# _tooling/scripts/kb-jobs-reap.timer
[Unit]
Description=Run kb-jobs-reap every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
